# example/supervisor/prometheusalert.ini
[program:PrometheusAlert]
command=/opt/prometheus/prometheusalert/PrometheusAlert
directory=/opt/prometheus/prometheusalert
stdout_logfile=/var/log/prometheus/PrometheusAlert.log
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
user=prometheus
